# crontab
* * * * * cd /var/www/html/mm-stock && /usr/bin/php artisan schedule:run >> /var/www/html/mm-stock/cron_output.log 2>&1
